# 2.10/files/services.ini
[supervisord]
nodaemon=false

[program:memcached]
command=/usr/bin/memcached -d -u memcached
startsecs=3
stopwaitsecs = 3

[program:rsyslog]
command=/usr/sbin/rsyslogd -n -d -c 2
startsecs = 5
stopwaitsecs = 5

[program:rsync]
command=/usr/bin/rsync --daemon --no-detach -v
startsecs = 5
stopwaitsecs = 5
